# home/private_dot_config/systemd/user/koi-drift-sweep.service
[Unit]
Description=KOI federation drift sweep
After=personal-koi.service

[Service]
Type=oneshot
ExecStart=%h/miniforge3/bin/python3 %h/.claude/local/scripts/koi_drift_sweep.py
StandardOutput=journal
StandardError=journal
