# home/private_dot_config/systemd/user/koi-canary.timer
[Unit]
Description=KOI canary daily timer

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
