# home/private_dot_config/systemd/user/koi-dlq-archive.timer
[Unit]
Description=KOI DLQ archive monthly timer

[Timer]
OnCalendar=*-*-01 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
