# home/private_dot_config/systemd/user/legion-hash-backfill.timer
[Unit]
Description=Run recording hash backfill nightly at 03:30 (P1.5)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
