# home/private_dot_config/systemd/user/legion-transcribe.service
[Unit]
Description=Legion Recording Transcription (batch)

[Service]
Type=oneshot
ExecStart=/bin/bash %h/.claude/local/scripts/transcribe-and-register.sh --batch-size 50
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
