# home/private_dot_config/systemd/user/legion-transcribe.service
[Unit]
Description=Legion Recording Transcription (with voice stagger)
# Whisper large-v3-turbo loads 2-4 GB resident per process. Caps below
# prevent OOM cascades that take down Spotify / Letta / etc. Two whisper
# instances ran concurrently 2026-04-27 17:32 (this batch + a one-off
# recover script) — kernel killed Spotify + the recover python.

[Service]
Type=oneshot
ExecStart=/bin/bash %h/.claude/local/scripts/transcribe-with-stagger.sh
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
MemoryHigh=6G
# Bumped 6G→9G after Apr 28 08:55 OOM mid-batch on 50-file run. faster-whisper
# peak ~4G + buffer; 9G gives headroom without encroaching on system RAM (32G total).
MemoryMax=9G
OOMScoreAdjust=500
