# home/private_dot_config/systemd/user/legion-diarize.service
[Unit]
Description=Legion Diarization Queue (with voice stagger — P1.1)
After=legion-transcribe.service
# OnFailure= triggers voice.target when the unit fails (e.g. kernel OOM
# SIGKILL bypasses bash trap in diarize-with-stagger.sh). Closes the
# stranded-voice gap (#321) at the systemd level — the bash trap covers
# normal/signal exits; this covers SIGKILL.
OnFailure=voice.target

[Service]
Type=oneshot
ExecStart=/bin/bash %h/.claude/local/scripts/diarize-with-stagger.sh
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
# WhisperX needs GPU memory — same limits as transcription
CPUQuota=100%
# Bumped 8G→12G after Apr 28 07:04 OOM on a long Mar-2021 .mkv. pyannote +
# WhisperX peak ~7-8G on multi-hour audio; 12G gives headroom without
# encroaching on system RAM (32G total).
MemoryMax=12G
Nice=10
