# home/private_dot_config/systemd/user/legion-pageindex-refresh.timer
[Unit]
Description=Daily PageIndex tree refresh

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
