# home/private_dot_config/systemd/user/legion-pipeline-watch.service
[Unit]
Description=Legion pipeline stall detector (P1.4)
After=default.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/uv run %h/.claude/local/scripts/pipeline-watch.py
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
MemoryMax=200M
Nice=15
