# home/private_dot_config/systemd/user/quartz-threads-adapter.timer
[Unit]
Description=Run Quartz threads adapter every 60 minutes (offset from skills adapter)

[Timer]
OnCalendar=*:34
Persistent=true
Unit=quartz-threads-adapter.service

[Install]
WantedBy=timers.target
