# home/private_dot_config/systemd/user/quartz-daily-digest.timer
[Unit]
Description=Nightly Quartz DailyDigest backfill (prior day) at 00:30

[Timer]
OnCalendar=*-*-* 00:30
Persistent=true
Unit=quartz-daily-digest.service

[Install]
WantedBy=timers.target
