# home/private_dot_config/systemd/user/quartz-baseline.timer
[Unit]
Description=Daily Quartz Phase 1 baseline at 23:50

[Timer]
OnCalendar=*-*-* 23:50:00
Persistent=true
Unit=quartz-baseline.service

[Install]
WantedBy=timers.target
