# home/private_dot_config/systemd/user/quartz-sessions-adapter.service
[Unit]
Description=Quartz sessions adapter — emit per-session pages from sessions/*.jsonl
After=quartz-rhythm.service

[Service]
Type=oneshot
WorkingDirectory=/home/shawn/.claude/plugins/local/legion-plugins/plugins/claude-quartz
ExecStart=%h/.local/bin/uv run scripts/adapters/sessions_to_markdown.py
StandardOutput=append:/home/shawn/.claude/local/quartz/sessions-adapter.log
StandardError=append:/home/shawn/.claude/local/quartz/sessions-adapter.log
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
MemoryMax=2G
TimeoutStartSec=10min
