# home/private_dot_config/systemd/user/quartz-session-replay.timer
[Unit]
Description=Run Quartz session-replay adapter every hour (after sessions adapter at *:30)

[Timer]
OnCalendar=*:36
Persistent=true
Unit=quartz-session-replay.service

[Install]
WantedBy=timers.target
