# home/private_dot_config/systemd/user/quartz-spectral-gap.service
[Unit]
Description=Quartz spectral-gap KPI emitter (Phase 3 D8)
After=default.target hippo-graph.service

[Service]
Type=oneshot
WorkingDirectory=%h/.claude/plugins/local/legion-plugins/plugins/claude-quartz
ExecStart=/usr/bin/env uv run --with falkordb --with scipy --with numpy --with pyyaml python %h/.claude/plugins/local/legion-plugins/plugins/claude-quartz/scripts/spectral_gap_emit.py
StandardOutput=append:%h/.claude/local/quartz/spectral-gap.log
StandardError=append:%h/.claude/local/quartz/spectral-gap.log
# Eigenvalue computation is single-threaded scipy ARPACK, ~9 min on 60K nodes.
TimeoutStartSec=20min
# Memory: sparse 60K × 60K Laplacian + ARPACK Krylov subspace ~ 1.5 GB peak.
MemoryHigh=2G
MemoryMax=3G
MemorySwapMax=0
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
