# home/private_dot_config/systemd/user/quartz-sessions-adapter.timer
[Unit]
Description=Run Quartz sessions adapter every 60 minutes (offset from internal build)

[Timer]
OnCalendar=*:30
Persistent=true
Unit=quartz-sessions-adapter.service

[Install]
WantedBy=timers.target
