# home/private_dot_config/systemd/user/quartz-threads-adapter.service
[Unit]
Description=Quartz threads adapter — emit curated message threads (default-deny allowlist)
After=quartz-rhythm.service

[Service]
Type=oneshot
WorkingDirectory=/home/shawn/.claude/plugins/cache/legion-plugins/claude-quartz/0.1.0
ExecStart=/usr/bin/uv run scripts/adapters/threads_to_markdown.py
StandardOutput=append:/home/shawn/.claude/local/quartz/threads-adapter.log
StandardError=append:/home/shawn/.claude/local/quartz/threads-adapter.log
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
MemoryMax=2G
TimeoutStartSec=5min
