# home/private_dot_config/systemd/user/quartz-session-replay.service
[Unit]
Description=Quartz session-replay adapter — Phase 7 D1 prompt-sequence JSON sidecars
After=quartz-sessions-adapter.service

[Service]
Type=oneshot
WorkingDirectory=/home/shawn/.claude/plugins/cache/legion-plugins/claude-quartz/0.1.0
ExecStart=/usr/bin/uv run scripts/adapters/session_replay_emit.py
StandardOutput=append:/home/shawn/.claude/local/quartz/session-replay-adapter.log
StandardError=append:/home/shawn/.claude/local/quartz/session-replay-adapter.log
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
MemoryMax=1G
TimeoutStartSec=10min
