# home/private_dot_config/systemd/user/quartz-temporal-slice.timer
[Unit]
Description=Run Quartz TemporalSlice emitter every hour (offset from build + sessions adapters)

[Timer]
OnCalendar=*:38
Persistent=true
Unit=quartz-temporal-slice.service

[Install]
WantedBy=timers.target
