# home/private_dot_config/systemd/user/quartz-rhythm.timer
[Unit]
Description=Run Quartz full build every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=quartz-rhythm.service

[Install]
WantedBy=timers.target
