# home/private_dot_config/systemd/user/rhythm-comms-digest.timer
[Unit]
Description=Communications Digest (daily, weekdays)

[Timer]
OnCalendar=Mon..Fri *-*-* 20:00:00 America/Vancouver
Persistent=true
Unit=rhythm@comms-digest.service

[Install]
WantedBy=timers.target
