# home/private_dot_config/systemd/user/rhythm-capacity-review.timer
[Unit]
Description=Nightly Capacity Review (11pm)

[Timer]
OnCalendar=*-*-* 23:00:00 America/Vancouver
Persistent=true
Unit=rhythm@capacity-review.service

[Install]
WantedBy=timers.target
