# home/private_dot_config/systemd/user/rhythm-evening-brief.timer
[Unit]
Description=Evening Brief (5pm)

[Timer]
OnCalendar=*-*-* 17:00:00 America/Vancouver
Persistent=true
Unit=rhythm@evening-brief.service

[Install]
WantedBy=timers.target
