# home/private_dot_config/systemd/user/rhythm-midday-brief.timer
[Unit]
Description=Midday Brief (noon)

[Timer]
OnCalendar=*-*-* 12:00:00 America/Vancouver
Persistent=true
Unit=rhythm@midday-brief.service

[Install]
WantedBy=timers.target
