# home/private_dot_config/systemd/user/venture-lead@.service
[Unit]
Description=Venture-lead persona wake — %i
Documentation=file://%h/.claude/local/research/2026/04/24/persona-mega-project/DESIGN.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.claude/local/scripts/venture-lead-driver.sh %i
TimeoutStartSec=300
Nice=5
StandardOutput=journal
StandardError=journal
Environment="HOME=%h"

[Install]
WantedBy=default.target
