# home/private_dot_config/systemd/user/venture-lead@bondometer.timer
[Unit]
Description=Bondometer (BCRG venture-lead) wake schedule
Requires=venture-lead@bondometer.service

[Timer]
# Mon/Wed/Fri 05:30 — 3x/week cadence (revised 2026-04-27 per Shawn)
OnCalendar=Mon,Wed,Fri *-*-* 05:30:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
