# home/private_dot_config/systemd/user/voice.target
[Unit]
Description=Voice Services (TTS + Arbiter)
Wants=voice-tts.service voice-arbiter.service
After=pipewire.service

[Install]
WantedBy=graphical-session.target
