# home/private_dot_config/systemd/user/voice-stt.service
[Unit]
Description=claude-voice STT daemon (wake word + transcription)
After=voice-queue.service pipewire.service
Requires=pipewire.service
PartOf=voice.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/home/shawn/.local/share/stt-env/bin/python3 /home/shawn/.claude/plugins/local/legion-plugins/plugins/claude-voice/scripts/stt_daemon.py
ExecStopPost=/bin/rm -f /home/shawn/.claude/local/voice/stt-active
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
KillSignal=SIGTERM
KillMode=process
Environment=PYTHONUNBUFFERED=1
# Route default mic capture through PipeWire WebRTC AEC source (task-398 P0).
# Subtracts Legion's own TTS (rendered to claude-voice-sink) from the mic
# input before wake-word detection sees it. Prevents recursive self-trigger.
# Verify pactl shows: voice-aec-source RUNNING. Override per-process via
# VOICE_AEC_DEVICE if PortAudio surfaces the node by name on this host.
Environment=PULSE_SOURCE=voice-aec-source

[Install]
WantedBy=voice.target
