# home/private_dot_config/systemd/user/voice-health.timer
[Unit]
Description=claude-voice health check (every 5 minutes)

[Timer]
OnBootSec=60
OnUnitActiveSec=300

[Install]
WantedBy=voice.target
