# home/private_dot_config/systemd/user/voice-queue-restart.service
[Unit]
Description=Restart queue daemon on code/config change

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart voice-queue.service
