# home/private_dot_config/systemd/user/voice-stt-proxy.service
[Unit]
Description=claude-voice STT WebSocket proxy (local transcription via faster-whisper)
After=pipewire.service
Wants=pipewire.service
PartOf=voice.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/home/shawn/.local/share/stt-env/bin/python3 \
    /home/shawn/.claude/plugins/local/legion-plugins/plugins/claude-voice/scripts/stt_proxy.py
ExecStopPost=/bin/rm -f /home/shawn/.claude/local/voice/stt-active
Restart=on-failure
RestartSec=5
TimeoutStartSec=10
KillSignal=SIGTERM
KillMode=mixed
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=4
Environment=CT2_CUDA_ALLOCATOR=cuda_malloc_async

# Memory guardrails
RuntimeMaxSec=43200
MemoryMax=6G
MemoryHigh=4G
MemorySwapMax=0

[Install]
WantedBy=voice.target
