# home/private_dot_config/systemd/user/voice.target
[Unit]
Description=Voice Services (TTS + STT + STT-proxy + Arbiter)
Wants=voice-tts.service voice-stt.service voice-stt-proxy.service voice-arbiter.service
After=pipewire.service

[Install]
WantedBy=graphical-session.target
