# home/private_dot_config/systemd/user/voice-tts-watch.path
[Unit]
Description=Watch for TTS daemon code/config changes

[Path]
PathChanged=/home/shawn/.claude/local/voice/config.yaml
PathChanged=/home/shawn/.claude/plugins/local/legion-plugins/plugins/claude-voice/scripts/tts_daemon.py
Unit=voice-tts-restart.service

[Install]
WantedBy=default.target
