# home/private_dot_config/systemd/user/voice-tts-restart.service
[Unit]
Description=Restart TTS daemon on code/config change

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart voice-tts.service
